# trash/memleak/memleak.service
[Unit]
Description=Memory leak daemon
After=sysinit.target

[Service]
ExecStart=/usr/sbin/memleak

[Install]
WantedBy=multi-user.target
