# nohang/openrc/nohang-desktop.in
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="nohang-desktop daemon"
description="Sophisticated low memory handler"
command=${BINDIR}/nohang
command_args="--config ${CONFDIR}/nohang/nohang-desktop.conf"
pidfile="/var/run/nohang-desktop"
start_stop_daemon_args="--background --make-pidfile"
