# nohang.service.in
[Unit]
Description=Highly configurable OOM prevention daemon
Documentation=man:nohang(1) https://github.com/hakavlad/nohang

[Service]
ExecStart=:TARGET_BIN:/nohang --config :TARGET_CONF:/nohang/nohang.conf
Slice=nohang.slice
Restart=always
KillMode=mixed
MemoryMax=100M
TasksMax=50
Nice=-15
CPUSchedulingResetOnFork=true

[Install]
WantedBy=multi-user.target
