# nohang.service
[Unit]
Description=No hang daemon
After=sysinit.target
Documentation=man:nohang(1) https://github.com/hakavlad/nohang

[Service]
ExecStart=/usr/local/bin/nohang
Restart=always

[Install]
WantedBy=multi-user.target
