# nohang.service.in
[Unit]
Description=Highly configurable OOM prevention daemon
Documentation=man:nohang(1) https://github.com/hakavlad/nohang

[Service]
ExecStart=:TARGET_BIN:/nohang --config :TARGET_CONF:/nohang/nohang.conf
Restart=always
KillMode=mixed
OOMPolicy=continue
Nice=-15
CPUSchedulingResetOnFork=true
MemoryMin=50%

[Install]
WantedBy=multi-user.target
