# nohang/nohang.service.in
[Unit]
Description=Sophisticated low memory handler
Documentation=man:nohang(1) https://github.com/hakavlad/nohang
Conflicts=nohang-desktop.service
After=system.slice

[Service]
ExecStart=:TARGET_BIN:/nohang --config :TARGET_CONF:/nohang/nohang.conf
SyslogIdentifier=nohang

KillMode=mixed

Restart=always
RestartSec=0

Nice=-5
CPUSchedulingResetOnFork=true
RestrictRealtime=yes

TasksMax=25
MemoryMax=250M

# Restrict access to the file system
UMask=0027
ReadOnlyPaths=/
ReadWritePaths=/var/log
InaccessiblePaths=/home /root

# Capabilities whitelist:
#  CAP_KILL is required to send signals (SIGTERM and SIGKILL)
#  CAP_IPC_LOCK  is required to mlockall()
#  CAP_SYS_PTRACE are required to check /proc/[pid]/exe realpathes
#  CAP_DAC_READ_SEARCH is required to read /proc/[pid]/environ files
#  CAP_DAC_READ_SEARCH CAP_AUDIT_WRITE CAP_SETUID CAP_SETGID CAP_SYS_RESOURCE are required to send GUI notifications

CapabilityBoundingSet=CAP_KILL CAP_IPC_LOCK CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_AUDIT_WRITE CAP_SETUID CAP_SETGID CAP_SYS_RESOURCE
AmbientCapabilities=CAP_KILL CAP_IPC_LOCK CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_AUDIT_WRITE CAP_SETUID CAP_SETGID CAP_SYS_RESOURCE

# It breaks GUI notifications on oldstable distros (Debian 8, CentOS 7)
PrivateNetwork=true

LockPersonality=yes
RestrictNamespaces=yes
ProtectKernelModules=true
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
